# pi_setup/supervisor/silvia_celery.conf
; ==================================
;  celery worker supervisor configuration
; ==================================

; WORKER 1 -------------------------------------------------
; the name of your supervisord program
[program:silvia_worker1]
; Set full path to celery program if using virtualenv
command=/home/ubuntu/.virtualenvs/venv-silvia/bin/celery worker -A silvia --queues=celery --loglevel=INFO
; The directory to your Django project
directory=/home/ubuntu/silvia/silvia
; If supervisord is run as the root user, switch users to this UNIX user account
; before doing any processing.
user=ubuntu
; Supervisor will start as many instances of this program as named by numprocs
numprocs=1
; Put process stdout output in this file
stdout_logfile=/var/log/celery/silvia_worker.log
; Put process stderr output in this file
stderr_logfile=/var/log/celery/silvia_worker.log
; If true, this program will start automatically when supervisord is started
autostart=true
; May be one of false, unexpected, or true. If false, the process will never
; be autorestarted. If unexpected, the process will be restart when the program
; exits with an exit code that is not one of the exit codes associated with this
; process’ configuration (see exitcodes). If true, the process will be
; unconditionally restarted when it exits, without regard to its exit code.
autorestart=true
; The total number of seconds which the program needs to stay running after
; a startup to consider the start successful.
startsecs=10
; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600
; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; WORKER 2 -------------------------------------------------
[program:silvia_worker2]
command=/home/ubuntu/.virtualenvs/venv-silvia/bin/celery worker -A silvia --queues=comms --loglevel=INFO --concurrency 1
directory=/home/ubuntu/silvia/silvia
user=ubuntu
numprocs=1
stdout_logfile=/var/log/celery/silvia_worker.log
stderr_logfile=/var/log/celery/silvia_worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true

; BEAT -------------------------------------------------
[program:silvia_celerybeat]
command=/home/ubuntu/.virtualenvs/venv-silvia/bin/celery beat -A silvia --scheduler django_celery_beat.schedulers:DatabaseScheduler --pidfile= --loglevel=INFO 
directory=/home/ubuntu/silvia/silvia
user=ubuntu
numprocs=1
stdout_logfile=/var/log/celery/silvia_beat.log
stderr_logfile=/var/log/celery/silvia_beat.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true

; GROUP -------------------------------------------------
; name of programs group
[group:silvia_workers]
; programs
programs=silvia_worker1,silvia_worker2,silvia_celerybeat
; if your broker is supervised, set its priority higher
; so it starts first
priority=998